# scripts/debian/postrm
#!/bin/sh
set -e

if [ "$1" = "purge" ] ; then
  update-rc.d dragent remove >/dev/null
fi
